# car.conf
[program:car]
command=/usr/bin/python3 /home/pi/car.py
autostart=true
